# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

i1	R1+	sbin/ini1
i2	R1+	sbin/ini2
i3	W1+	sbin/ini3

s1	F3+	sbin/srv1
s2	F3+	sbin/srv2
s3	F3+	sbin/srv3

	R789	!echo going to sleep
sleep	W7	!echo mem > /sys/power/state
	X789	!echo wakeup

# vim: ft=
