# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

E	i1	sbin/ini1
E	i2	sbin/ini2
W	i3	sbin/ini3

S	s1	sbin/srv1
S	s2	sbin/srv2
S	s3	sbin/srv3

W8	ss	!echo sleeping
W8~	sw	!echo wakeup

# vim: ft=
