# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

ini1	R1+	sbin/ini1
ini2	R1+	sbin/ini2
ini3	W1+	sbin/ini3

srv1	F3+	sbin/srv1
srv2	F3+	sbin/srv2
srv3	F3+	sbin/srv3

	R789	!echo going to sleep
sleep	W7	!echo mem > /sys/power/state
	X789	!echo wakeup

# vim: ft=
