# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

::once:sbin/ini1
::once:sbin/ini2
::wait:sbin/ini3

s1:::sbin/srv1
s2:::sbin/srv2
s3:a::sbin/srv3

:8:wait:!echo sleeping
:~8:wait:!echo wakeup
