# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

r	i1	sbin/ini1
r	i2	sbin/ini2
R	i3	sbin/ini3

s	s1	sbin/srv1
s	s2	sbin/srv2
s	s3	sbin/srv3

R8	ss	!echo sleeping
R8~	sw	!echo wakeup

# vim: ft=
