# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

i1	r	sbin/ini1
i2	r	sbin/ini2
i3	w	sbin/ini3

s1	-	sbin/srv1
s2	-	sbin/srv2
s3	-	sbin/srv3

ss	w8	!echo sleeping
sw	w8~	!echo wakeup

# vim: ft=
