# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

r	i1	sbin/ini1
r	i2	sbin/ini2
w	i3	sbin/ini3

-	s1	sbin/srv1
-	s2	sbin/srv2
-	s3	sbin/srv3

w8	ss	!echo sleeping
w8~	sw	!echo wakeup

# vim: ft=
