# etc/inittab
# inittab for in-place testing

PATH=/bin
VAR=value

# System boot
mount	W1+	sbin/mount -o remount,rw /
hwc	R1+	sbin/hwclock
	W1+	sbin/net up

# System shutdown
	W0	sbin/net down
	W0	sbin/umount -a

# Respawning services (see also initdir/)
srv1	S3+	sbin/trap srv1 10
srv2	S3+	sbin/trp srv2 23
srv3	S3+	sbin/trap srv3

# Sleep modes
	R789	sbin/net down
sleep	W8	sbin/sleep 2
suspend	W9	sbin/sleep 5
	X789	sbin/net up

# vim: ft=conf
