# misc/buildroot/board/sninit/root.aarch64/etc/inittab
PATH=/bin:/sbin

proc	R	/bin/mount -t proc proc /proc
pts	R	/bin/mkdir -p /dev/pts
shm	R	/bin/mkdir -p /dev/shm
host	R	/bin/hostname -F /etc/hostname

syslog	L	/sbin/syslogd -n
klogd	S	/sbin/klogd -nc1
dhcp	L	/sbin/run -n /sbin/udhcpc -fRS eth0

ttyS	T	/sbin/getty -L ttyAMA0 115200 vt100

# vim: ft=conf:
