# misc/buildroot/board/sninit/root/etc/inittab
PATH=/bin:/sbin

proc	R	/bin/mount -t proc proc /proc
pts	R	/bin/mkdir -p /dev/pts
shm	R	/bin/mkdir -p /dev/shm
mount	W	/bin/mount -o remount,rw /
host	R	/bin/hostname -F /etc/hostname

syslog	L	/sbin/syslogd -n
klogd	S	/sbin/klogd -nc1
dhcp	L	/sbin/run -n /sbin/udhcpc -fRS eth0

tty1	T	/sbin/getty -L tty1 115200 linux
tty2	T	/sbin/getty -L tty2 115200 linux
tty3	T	/sbin/getty -L tty3 115200 linux

umount	X	/bin/umount -ar

# vim: ft=conf:
