# misc/buildroot/board/sninit/root/etc/inittab
PATH=/bin:/sbin

proc	R	/bin/mount -t proc proc /proc
pts	R	/bin/mkdir -p /dev/pts
shm	R	/bin/mkdir -p /dev/shm
mount	R	/bin/mount -o remount,rw /
host	R	/bin/hostname -F /etc/hostname

syslog	L	/sbin/syslogd -n -O /var/log/messages
klogd	S	/sbin/klogd -nc1
dhcp	L	/sbin/run -n /sbin/udhcpc -fRS eth0

tty1	T	/sbin/getty -L tty1 115200 vt100

umount	X	/bin/umount -ar

# vim: ft=conf:
